# packaging/deb/debian/prerm
#!/usr/bin/env bash
set -euo pipefail

case "${1:-}" in
  remove|deconfigure|upgrade|failed-upgrade)
    systemctl disable --now vulhunter.service || true
    ;;
esac

exit 0
